# dst-caves.service
[Unit]
Description=dst caves

[Service]
Restart=always
RestartSec=5s
User=dst2
WorkingDirectory=/home/dst2/server_dst/bin
ExecStart=sh start_caves.sh

[Install]
WantedBy=multi-user.target

